# DublinBikes/dublin_bikes/historic_data.cron
*/5 * * * * /usr/bin/python /home/ec2-user/Group8Project/dublin_bikes/historic_data.py