# DublinBikes/dublin_bikes/realtime.cron
#This file runs every minute to collect real time bike and weather info.

*/1 * * * * /usr/bin/python /home/ec2-user/Group8Project/dublin_bikes/update.py
*/1 * * * * /usr/bin/python /home/ec2-user/Group8Project/dublin_bikes/weatherscraper.py